# agents/systemd/eg-slot@.service
[Unit]
Description=EG Slot Agent %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/eg/agents/slot
ExecStart=/usr/bin/python3 /opt/eg/agents/slot/slot_agent.py /opt/eg/agents/slot/configs/%i.yaml
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
